# crontab
# crontab
# 0 3 * * * /bin/sh /backup.sh >> /var/log/backup.log 2>&1
* * * * * /bin/sh /backup.sh >> /var/log/backup.log 2>&1
